# webhooks/github_pull_request/crontab
SHELL=/bin/bash
*/2 * * * * sh ~/webhook/check_pr.sh
